# supervisord-php5-fpm.conf
[program:php]
command=/start-php5-fpm.sh
autostart=true
autorestart=true
