# supervisord-nginx.conf
[program:nginx]
command=/start-nginx.sh
numprocs=1
autostart=true
autorestart=true
